# sampleruby.service.j2
[Unit]
Description=Simple service to run ruby script
After=network.target

[Service]
Type=simple
ExecStart=/home/vagrant/sampleruby
PIDFile=/var/run/sampleruby

[Install]
WantedBy=default.target
